# image/fs-config/lib/udev/rules.d/99-usb-leds.rules
ACTION=="add|remove", SUBSYSTEM=="usb", RUN+="usb-leds.sh"
